# crontab
@reboot nginx &>/var/log/nginx.cron.log &
*/15 * * * * /app/scanner/main.py /app/web/albums /app/web/cache &>/var/log/scanner.cron.log
