# crontab
@reboot root nginx &>/var/log/nginx.cron.log &
*/15 * * * * root /app/scanner/main.py /app/web/albums /app/web/cache &>/var/log/scanner.cron.log
